# filebeat/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/var/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=true                ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024

[program:daemon]
user=root
command=filebeat --path.config /usr/share/filebeat | python console2rabbit.py
stderr_logfile=/var/log/%(program_name)s.log
autostart=true
autorestart=true